# builds/any/rootfs/stretch/common/overlay/etc/systemd/system/networking.service.d/switchdev-online.conf
[Unit]
After=local-fs.target network-pre.target
Before=shutdown.target network.target network-online.target
